# philips_hue_data_collector.timer
[Unit]
Description=Timer for Philips Hue Data Collector

[Timer]
OnCalendar=*:*:00
Persistent=true

[Install]
WantedBy=timers.target